# gitScript.service
[Unit]
Description=Git Script
[Service]
Type=oneshot
ExecStart=/home/mankirat/Unix_Final_Project/gitScript.sh
